# deploy/riyad-discord-bot.service.example
[Unit]
Description=Riyad Discord Bot
After=network.target

[Service]
Type=simple
User=opc
WorkingDirectory=/opt/riyad-discord-bot
ExecStart=/usr/bin/node src/bootstrap.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
